# servers/benolilab/vmdocker/docker-compose/monitorance/volume-backup-wrapper/backup.cron
# Run at the specified time in BACKUP_SCHEDULE - supercronic will execute the schedule - see README at https://github.com/aptible/supercronic
# Note: job labels help the preview-schedule.sh view each job with a label to help identify/verify it
# job: backup-job
$BACKUP_SCHEDULE /app/scripts/run-backup.sh
# job: healthcheck-heartbeat
0 0 * * * /app/scripts/healthcheck-heartbeat.sh
